# etc/profile.d/environment.sh
export PUID=
export PGID=
export DOMAIN=
export APPDATA=
export CF_EMAIL=
export CF_API_KEY=
export NAMESERVER_1=
export NAMESERVER_2=
export ACME_EMAIL=
export LOGIN_OAUTH_SECRET=
export LOGIN_OAUTH_CLIENT_ID=
export LOGIN_OAUTH_CLIENT_SECRET=
export LOGIN_OAUTH_WHITELIST=
export MNT_HDD=
export TORRENT_PORT=
export MINECRAFTBED_PORT=
export PLEX_CLAIM=
